# com.set.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique job label — must match the filename -->
    <key>Label</key>
    <string>com.set.monitor</string>

    <!-- Script to run (sources set_env.sh for credentials) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/_pporpin_/Documents/Claude/Projects/Stock/set_run.sh</string>
    </array>

    <!-- Run every 30 minutes (1800 seconds) -->
    <!-- The monitor itself exits immediately if market is closed -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Don't run immediately when loaded — wait for first interval -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Log output to /tmp for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/set_monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/set_monitor_err.log</string>

    <!-- Restart automatically if the script crashes (not just exits cleanly) -->
    <!-- Note: a clean exit (market closed) does NOT trigger this -->
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

</dict>
</plist>
